# streamlit-optimization.service
[Unit]
Description=Streamlit BellaTrix-v1 App
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/BellaTrix-v1
Environment=PATH=/home/ubuntu/BellaTrix-v1/.venv/bin:/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin
# Use wrapper script that loads GitHub Secrets
ExecStart=/home/ubuntu/BellaTrix-v1/scripts/start-streamlit-with-secrets.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

